# systemd/remind.timer
[Unit]
Description=Morning and evening reminders

[Timer]
OnCalendar=*-*-* 09:00:00
OnCalendar=*-*-* 18:00:00
Persistent=true

[Install]
WantedBy=timers.target
